# jenkins/ansible/rc.local
#!/bin/sh -e

# setup network namespaces
/usr/local/bin/netns-setup

# resize /run
mount -oremount /run

exit 0

