# repo/crontab
# do not edit, this file is stored in git

MAILTO=cb@df7cb.de

# collect statistics about the repository
0 0 * * *  /srv/apt/apt.postgresql.org/repo/apt-stats

# update dists tarball for the debcheck jenkins jobs
13 * * * *  tar -C /srv/apt/pub/repos/apt/ -cf /srv/apt/web/dists.tar.gz dists/

# move morgue files into pool directories
16 * * * *  /srv/apt/apt.postgresql.org/repo/bin/morguefiles-to-pool

# remove snapshot builds from morgue to save some space
19 0 * * *  find /srv/apt/repo/morgue/p/postgresql-1[01]/ -name '*~*~*' -mtime +14 | xargs -r rm
